# ops/systemd/tifa-assistant-web.service
[Unit]
Description=Tifa Assistant Framework Web App
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=${TIFA_SERVICE_USER}
WorkingDirectory=${TIFA_REPO_ROOT}
Environment=NODE_ENV=production
Environment=HOST=0.0.0.0
Environment=PORT=3100
Environment=TRADEVIBE_RUNTIME_DIR=${TIFA_REPO_ROOT}/runtime
Environment=TRADEVIBE_TIMEZONE=Asia/Ho_Chi_Minh
EnvironmentFile=-${TIFA_REPO_ROOT}/.env
ExecStart=/bin/bash ${TIFA_REPO_ROOT}/start.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

