# ops/systemd/tradevibe-insight.service
[Unit]
Description=TradeVibe Daily Insight Pipeline

[Service]
Type=oneshot
User=nexus
WorkingDirectory=/home/nexus/projects/tradevibe-org
Environment="TRADEVIBE_RUNTIME_DIR=runtime"
ExecStart=/bin/bash /home/nexus/projects/tradevibe-org/insight_engine/tradevibe_runner.sh

[Install]
WantedBy=multi-user.target
