# ops/systemd/tradevibe-insight.timer
[Unit]
Description=Run TradeVibe insight pipeline every day at 08:00

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
Unit=tradevibe-insight.service

[Install]
WantedBy=timers.target
