# deploy/systemd/tifa-tts-worker.service
[Unit]
Description=Tifa Assistant Framework TTS Worker
After=network.target

[Service]
Type=simple
User=nexus
WorkingDirectory=/home/nexus/projects/tifa-assistant-framework
Environment=NODE_ENV=production
Environment=PATH=/home/nexus/.nvm/versions/node/v22.18.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=/home/nexus/projects/tifa-assistant-framework/.env
ExecStart=/home/nexus/.nvm/versions/node/v22.18.0/bin/npm run tts:worker
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/home/nexus/projects/tifa-assistant-framework

[Install]
WantedBy=multi-user.target
